# smart_energy_controller/rootfs/etc/services.d/smart-energy/finish
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Finish script for Smart Energy Controller
# ==============================================================================

bashio::log.info "Smart Energy Controller stopped."
